# maahi/com.meet.maahi.consolidate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Memory consolidator — launchd plist
  Runs at 03:00 daily and turns yesterday's conversation JSONLs into durable
  facts + preferences inside Maahi's vault memory.

  Install:
    cp com.meet.maahi.consolidate.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.meet.maahi.consolidate.plist
  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.meet.maahi.consolidate.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meet.maahi.consolidate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "/Users/themeetpatel/My brain/maahi" &amp;&amp; ./.venv/bin/python -m maahi.consolidate</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/themeetpatel/My brain/maahi/logs/consolidate.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/themeetpatel/My brain/maahi/logs/consolidate.err.log</string>
</dict>
</plist>
